# deploy/systemd/ambient-bridge.service
[Unit]
Description=GENesis-Voice ambient bridge (capture -> Silero VAD -> Zipformer STT -> ambient.db)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Point WorkingDirectory at the parent of the deployed ambient_bridge/ package
# so `python -m ambient_bridge.server` resolves.
WorkingDirectory=%h/genesis-voice/bridges
# Use your venv's python (the bridge needs sherpa-onnx, onnxruntime, soxr, websockets).
ExecStart=%h/ambient-venv/bin/python -m ambient_bridge.server
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Cap glibc malloc arenas. sherpa/onnxruntime allocate large, variable per-inference buffers on
# the to_thread STT pool (and the diar spawn child, which inherits this env); glibc gives each
# pool thread its own arena and fragmented blocks never return to the OS, so RSS creeps for days.
# Capping to 2 arenas bounds it (validated on-edge: RSS plateaus vs an unbounded ~200 MB/hr climb).
Environment=MALLOC_ARENA_MAX=2
# Override model paths / tuning here if non-default, e.g.:
# Environment=AMBIENT_INPUT_SR=24000
# Environment=AMBIENT_VAD_MIN_SILENCE=0.4
# The RESIDUAL RSS ratchet (survives MALLOC_ARENA_MAX=2) is onnxruntime's BFC arena growing
# under variable-length utterances. Opt out for the variable-shape sessions (README "Memory";
# ~4.5% RTF cost, flat RSS — measured on-edge) and/or cap the diar child's RSS (recycles the
# subprocess between windows above the ceiling):
# Environment=AMBIENT_ORT_ARENA_OFF=1
# Environment=AMBIENT_DIAR_RSS_CEILING_MB=1400

[Install]
WantedBy=default.target
