# bridges/ambient_bridge/systemd/ambient-bridge.service
[Unit]
Description=Genesis Ambient Bridge (capture -> Silero VAD -> Zipformer STT -> ambient.db)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
# %h = the service user's home; the deployed ambient_bridge/ package sits there.
ExecStart=%h/bench/bin/python -m ambient_bridge.server
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Override model paths / tuning here if non-default, e.g.:
# Environment=AMBIENT_INPUT_SR=24000
# Environment=AMBIENT_VAD_MIN_SILENCE=0.4

[Install]
WantedBy=default.target
